# ZimbraServer/src/bin/zimbra
#!/bin/bash
#
# Init file for zimbra mail
#
# chkconfig: 345 99 01
# description: Zimbra mail service
#


command()
{
	su - zimbra -c "zmcontrol $1"
}

case "$1" in
	restart)
		command shutdown
		command startup
		RETVAL=$?
		;;
	start)
		command startup
		RETVAL=$?
		;;
	stop)
		command shutdown
		RETVAL=$?
		;;
	reload|status)
		command $1
		RETVAL=$?
		;;
	*)
		echo $"Usage: $0 {start|stop|restart|reload|status}"
		RETVAL=1
esac
exit $RETVAL
